# curator4-871559621.conf
[program:curator4-871559621]
command=python3 /home/curator4/curator4-871559621/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-871559621/test.err.log
stdout_logfile=/home/curator4/curator4-871559621/test.out.log